# roles/generic/systemd-status-mail/templates/status-mail@.service.jinja2
[Unit]
Description=Status mail for %i

[Service]
Type=oneshot
ExecStart={{ systemd_status_mail_script }} %i
User=nobody
