# playbooks/roles/generic/smb-mount/templates/smb-mount-fixup.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{ smb_mount_fixup_label }}</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ smb_mount_fixup_script }}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{ smb_mount_fixup_log }}</string>
    <key>StandardErrorPath</key>
    <string>{{ smb_mount_fixup_log }}</string>
</dict>
</plist>
